# build/audio-player_1.0.4/DEBIAN/prerm
#!/bin/bash
set -e

echo "Stopping audio-player service..."
systemctl stop audio-player.service || true
systemctl disable audio-player.service || true

exit 0
