# build/audio-player_1.0.0/DEBIAN/prerm
#!/bin/bash
set -e

# Stop and disable service
systemctl stop audio-player.service || true
systemctl disable audio-player.service || true

exit 0
